# systemd/hermes-smoke.timer
[Unit]
Description=Run Hermes per-profile smoke daily (08:40)

[Timer]
OnCalendar=*-*-* 08:40:00
Persistent=true
AccuracySec=5m
Unit=hermes-smoke.service

[Install]
WantedBy=timers.target
