# systemd/hermes-lessons.service
[Unit]
Description=Hermes weekly graduation-lessons (model-driven reflection -> Telegram)
After=network-online.target

[Service]
Type=oneshot
Environment=LANG=C.UTF-8
Environment=LC_ALL=C.UTF-8
Environment=PYTHONIOENCODING=utf-8
ExecStart=/root/hermes/runtime/bin/hermes_lessons_alert.sh
