# systemd/memory-index.timer
[Unit]
Description=Weekly rebuild of orchestrator local memory index

[Timer]
OnCalendar=Sun 09:40
Persistent=true

[Install]
WantedBy=timers.target
